# systemd/auditslip-dashboard.service
[Unit]
Description=Auditslip Dashboard
After=network-online.target auditslip-bot.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/root/projects/auditslip
EnvironmentFile=/etc/auditslip/auditslip.env
ExecStart=/usr/bin/python3 /root/projects/auditslip/auditslip_dashboard.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
